# files/etc/cron.d/tarsnap
# Run backup script
* */12 * * * root /usr/local/bin/tarsnap.sh > /dev/null
