# files/etc/cron.d/tarsnap
# Run backup script
*/5 * * * * root /usr/local/bin/tarsnap.sh > /dev/null
